# hw07/uwsgi-config/uwsgi.service.d/changes.conf
[Unit]
After=network-online.target

[Service]
ExecStart=
ExecStart=/usr/sbin/uwsgi --ini /etc/uwsgi.ini:former
