# hw12/systemd/zabbix-server.service.d/changes.conf
[Install]
WantedBy=multi-user.target monitoring.target

